# debian/rfidtime-bifferboard.conffiles
/etc/init.d/rfidtime-bifferboard
/etc/cron.d/rfidtime-bifferboard
/etc/modprobe.d/rfidtime-bifferboard.conf

