# debian/rfidtime-alix.conffiles
/etc/init.d/rfidtime-alix
/etc/cron.d/rfidtime-alix

