# debian/rfidtime-bifferboard.cron.d
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=richi@paraeasy.ch
HOME=/tmp

* * * * * /etc/init.d/rfidtime-bifferboard start
