# debian/rfidtime-bifferboard.prerm
#!/bin/sh

set -e

case "$1" in
	remove|upgrade|remove-in-favour|deconfigure-in-favour)
		#
		# do here pre removes
		#
		update-rc.d -f rfidtime-bifferboard remove

		;;
	failed-upgrade)
		;;
esac

exit 0
